# myapp.service
[Unit]
Description=My Go App
After=network.target

[Service]
Type=simple
User=myapp
ExecStartPre=/usr/sbin/id -u myapp || /usr/sbin/useradd -r -s /bin/false myapp
ExecStart=/opt/myapp/main
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
